# docker/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
username = dummy
password = dummy

[inet_http_server]
port=127.0.0.1:9001
username = dummy
password = dummy

[supervisord]
nodaemon=true
loglevel=warn

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
username = dummy
password = dummy


# [eventlistener:flush_signalbus_trigger]
# command=%(ENV_APP_ROOT_DIR)s/trigger_supervisor_process.py signalbus_flush 60
# directory=%(ENV_APP_ROOT_DIR)s
# events=TICK_5

# [program:signalbus_flush]
# command=flask signalbus flush -w 10
# directory=%(ENV_APP_ROOT_DIR)s
# autostart=false
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes = 0
# redirect_stderr=true
# startsecs=0
# exitcodes=0
# autorestart=false


[program:flush_rejected_transfers]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh flush_rejected_transfers
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:flush_prepared_transfers]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh flush_prepared_transfers
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:flush_finalized_transfers]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh flush_finalized_transfers
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:flush_account_transfers]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh flush_account_transfers
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:flush_account_updates]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh flush_account_updates
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:flush_account_purges]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh flush_account_purges
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:flush_rejected_configs]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh flush_rejected_configs
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:flush_pending_balance_changes]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh flush_pending_balance_changes
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:process_balance_changes]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh process_balance_changes
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:process_transfer_requests]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh process_transfer_requests
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:process_finalization_requests]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh process_finalization_requests
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_accounts]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_accounts
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_registered_balance_changes]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_registered_balance_changes
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:scan_prepared_transfers]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh scan_prepared_transfers
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:dramatiq_protocol]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh protocol
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:process_chores]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh process_chores
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000


[program:webserver]
command=%(ENV_APP_ROOT_DIR)s/entrypoint.sh webserver
directory=%(ENV_APP_ROOT_DIR)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=30
startretries=1000000
